# resources/supervisor/cosmian_kms.ini
[program:cosmian_kms]
command=cosmian_kms
environment=COSMIAN_KMS_CONF=/var/lib/cosmian_vm/data/app.conf
directory=/usr/sbin
autostart=false
autorestart=true
user=root
stderr_logfile=/var/log/cosmian_kms.err.log
stdout_logfile=/var/log/cosmian_kms.out.log
