# resources/supervisor/cosmian_kms.ini
[program:cosmian_kms]
command=cosmian_kms
directory=/usr/sbin
autostart=false
autorestart=true
user=root
stderr_logfile=/var/log/cosmian_kms.err.log
stdout_logfile=/var/log/cosmian_kms.out.log
