# pkg/deb/postrm
#!/bin/sh

has_systemd() {
  command -v systemctl >/dev/null 2>&1 && [ -d /run/systemd/system ]
}

if ! has_systemd; then
  exit 0
fi

systemctl daemon-reload >/dev/null 2>&1 || true

# Check if the service is active (running)
if systemctl is-active --quiet cosmian_kms.service; then
  echo "Service is running. Stopping the service..."
  systemctl stop cosmian_kms.service >/dev/null 2>&1 || true

  # Verify if the service was stopped
  if systemctl is-active --quiet cosmian_kms.service; then
    echo "Failed to stop the service."
  else
    echo "Service stopped successfully. Disabling it"
    if systemctl list-unit-files | grep -q '^cosmian_kms.service'; then
      echo "Service exists. Disabling the service..."
      systemctl disable cosmian_kms.service >/dev/null 2>&1 || true
    fi
  fi
else
  echo "Service is not running. Disabling it"
  if systemctl list-unit-files | grep -q '^cosmian_kms.service'; then
    echo "Service exists. Disabling the service..."
    systemctl disable cosmian_kms.service >/dev/null 2>&1 || true
  fi
fi

systemctl daemon-reload >/dev/null 2>&1 || true

#DEBHELPER#
